# bambuddy/rootfs/etc/services.d/bambuddy/run
#!/usr/bin/with-contenv bashio

TZ=$(bashio::config 'timezone')
LOG_LEVEL=$(bashio::config 'log_level')

export TZ="${TZ:-Europe/Rome}"
export LOG_LEVEL="${LOG_LEVEL:-info}"

if bashio::config.has_value 'bind_address'; then
    BIND_ADDRESS=$(bashio::config 'bind_address')
    export BIND_ADDRESS="${BIND_ADDRESS}"
    bashio::log.info "Binding to: ${BIND_ADDRESS}"
fi

mkdir -p /data/bambuddy_data /data/bambuddy_logs

if [ ! -L /app/data ]; then
    rm -rf /app/data
    ln -s /data/bambuddy_data /app/data
fi
if [ ! -L /app/logs ]; then
    rm -rf /app/logs
    ln -s /data/bambuddy_logs /app/logs
fi

export PORT=8000
export PYTHONPATH=/app

bashio::log.info "Starting BamBuddy on port ${PORT}..."

cd /app
exec python3 -m uvicorn backend.app.main:app \
    --host 0.0.0.0 \
    --port "${PORT}" \
    --loop asyncio
