# bambuddy/rootfs/etc/services.d/bambuddy/run
#!/usr/bin/with-contenv bashio


LOG_LEVEL=$(bashio::config 'log_level')
export LOG_LEVEL="${LOG_LEVEL:-info}"

# Home Assistant integration.
# If ha_url/ha_token are not set, fall back to the Supervisor's internal
# Core API proxy (available because homeassistant_api: true in config.yaml
# makes the Supervisor inject SUPERVISOR_TOKEN into this container).
if bashio::config.has_value 'ha_url'; then
    HA_URL=$(bashio::config 'ha_url')
else
    HA_URL="http://supervisor/core"
fi
export HA_URL
bashio::log.info "Home Assistant URL: ${HA_URL}"

if bashio::config.has_value 'ha_token'; then
    HA_TOKEN=$(bashio::config 'ha_token')
else
    HA_TOKEN="${SUPERVISOR_TOKEN}"
fi
export HA_TOKEN

HA_INFO_RESPONSE=$(curl --silent --fail \
  --header "Authorization: Bearer ${HA_TOKEN}" \
  "${HA_URL%/core}/info" 2>/dev/null) || true

TZ=""
if [ -n "${HA_INFO_RESPONSE}" ]; then
    TZ=$(bashio::jq "${HA_INFO_RESPONSE}" '.data.timezone // empty')
fi
export TZ="${TZ:-UTC}"
bashio::log.info "Timezone: ${TZ}"

if bashio::config.has_value 'bind_address'; then
    BIND_ADDRESS=$(bashio::config 'bind_address')
    export BIND_ADDRESS="${BIND_ADDRESS}"
    bashio::log.info "Binding to: ${BIND_ADDRESS}"
fi

mkdir -p /data/bambuddy_data /data/bambuddy_logs

if [ ! -L /app/data ]; then
    rm -rf /app/data
    ln -s /data/bambuddy_data /app/data
fi
if [ ! -L /app/logs ]; then
    rm -rf /app/logs
    ln -s /data/bambuddy_logs /app/logs
fi

# Fix finish permissions at runtime
chmod +x /etc/services.d/bambuddy/finish 2>/dev/null || true

if bashio::config.has_value 'trusted_frame_origins'; then
    TRUSTED_FRAME_ORIGINS=$(bashio::jq "/data/options.json" '.trusted_frame_origins // [] | join(",")')
    if [ -n "${TRUSTED_FRAME_ORIGINS}" ]; then
        export TRUSTED_FRAME_ORIGINS
        bashio::log.info "Trusted frame origins: ${TRUSTED_FRAME_ORIGINS}"
    fi
fi

if bashio::config.has_value 'database_url'; then
    DATABASE_URL=$(bashio::config 'database_url')
    export DATABASE_URL="${DATABASE_URL}"
    bashio::log.info "Using external PostgreSQL database"
fi

# External storage (Share / Media). Mounted via config.yaml (share:rw, media:rw);
# these two toggles just control whether BamBuddy is allowed to register them.
EXTERNAL_ROOTS=""
if bashio::config.true 'enable_share'; then
    bashio::log.info "Enabling /share as external root"
    EXTERNAL_ROOTS="/share"
fi
if bashio::config.true 'enable_media'; then
    bashio::log.info "Enabling /media as external root"
    if [ -n "${EXTERNAL_ROOTS}" ]; then
        EXTERNAL_ROOTS="${EXTERNAL_ROOTS}:/media"
    else
        EXTERNAL_ROOTS="/media"
    fi
fi
if [ -n "${EXTERNAL_ROOTS}" ]; then
    export BAMBUDDY_EXTERNAL_ROOTS="${EXTERNAL_ROOTS}"
    bashio::log.info "External roots: ${EXTERNAL_ROOTS}"
fi

if bashio::config.true 'use_system_trust_store'; then
    CERTFILE=$(bashio::config 'certfile')
    CERT_SRC="/config/${CERTFILE}"
    if [ -n "${CERTFILE}" ] && [ -f "${CERT_SRC}" ]; then
        bashio::log.info "Installing CA certificate: ${CERTFILE}"
        cp "${CERT_SRC}" "/usr/local/share/ca-certificates/${CERTFILE}"
        update-ca-certificates
        export USE_SYSTEM_TRUST_STORE="true"
        bashio::log.info "System trust store enabled"
    else
        bashio::log.warning "use_system_trust_store is enabled but certificate file not found: ${CERT_SRC}"
        bashio::log.warning "Place your CA certificate (.crt) in the addon_configs folder and set certfile to the correct filename"
    fi
fi

export PORT=8000
export PYTHONPATH=/app
export PYTHONUNBUFFERED=1

bashio::log.info "Starting BamBuddy on port ${PORT}..."

cd /app
exec python3 -m uvicorn backend.app.main:app \
    --host 0.0.0.0 \
    --port "${PORT}" \
    --loop asyncio