# bambuddy/rootfs/etc/services.d/bambuddy/run
#!/bin/bash
set -e

OPTIONS_FILE="/data/options.json"

if [ -f "$OPTIONS_FILE" ]; then
    TZ=$(python3 -c "import json; d=json.load(open('$OPTIONS_FILE')); print(d.get('timezone','Europe/Rome'))")
    LOG_LEVEL=$(python3 -c "import json; d=json.load(open('$OPTIONS_FILE')); print(d.get('log_level','info'))")
    BIND_ADDRESS=$(python3 -c "import json; d=json.load(open('$OPTIONS_FILE')); print(d.get('bind_address',''))" 2>/dev/null || echo "")
fi

export TZ="${TZ:-Europe/Rome}"
export LOG_LEVEL="${LOG_LEVEL:-info}"

if [ -n "$BIND_ADDRESS" ]; then
    export BIND_ADDRESS
    echo "[bambuddy] Binding to address: ${BIND_ADDRESS}"
fi

# Persistent storage
mkdir -p /data/bambuddy_data /data/bambuddy_logs

if [ ! -L /app/data ]; then
    rm -rf /app/data
    ln -s /data/bambuddy_data /app/data
fi
if [ ! -L /app/logs ]; then
    rm -rf /app/logs
    ln -s /data/bambuddy_logs /app/logs
fi

export PORT=8000

echo "[bambuddy] Starting BamBuddy on port ${PORT}..."

# Bypass gosu by running as root (HA addon containers run as root anyway)
# Skip the entrypoint entirely and start uvicorn directly
exec python3 -m uvicorn backend.app.main:app \
    --host 0.0.0.0 \
    --port ${PORT} \
    --loop asyncio
